# examples/ycast.service.example
[Unit]
Description=YCast internet radio service
After=network.target

[Service]
User=root
ExecStart=/opt/ycast/ycast.py

[Install]
WantedBy=multi-user.target
